# cmd/downloader/crontab/open-weather-api-downloader-cron
*/30 * * * * cd /Users/pankaj/Documents/Work/Open-Weather-API/cmd/downloader && go run . >> /var/log/open-weather-api-downloader-cron.log 2>&1
